# charts-syncer/crontab
* */6 * * * charts-syncer sync --latest-version-only --config /tmp/charts-syncer.yaml
# Don't remove the empty line at the end of this file. It is required to run the cron job